# bin/update-google-api.cron
#! /bin/sh
cd `dirname $0`/../../..
/usr/bin/env php -f ./symfony opCalendar:update-google-api $1
